# scripts/com.agent-infra.hook-roi-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-infra.hook-roi-daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/alien/Projects/agent-infra &amp;&amp; uv run python3 scripts/hook-roi.py --days 1 > "artifacts/hook-roi/$(date +%Y-%m-%d).md" 2>&amp;1</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/alien/.claude/hook-roi-daily-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alien/.claude/hook-roi-daily-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/alien/.local/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/alien</string>
    </dict>
</dict>
</plist>
