# scripts/com.meta.propose-work-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meta.propose-work-daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/alien/Projects/meta &amp;&amp; uv run python3 scripts/propose-work.py --days 1 2>&amp;1 | tee /Users/alien/.claude/propose-work.log</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/alien/.claude/propose-work-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alien/.claude/propose-work-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/alien/.local/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/alien</string>
    </dict>
</dict>
</plist>
