# scripts/com.meta.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meta.orchestrator</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/alien/Projects/meta/scripts/orchestrator.py</string>
        <string>tick</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/Users/alien/.claude/orchestrator-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alien/.claude/orchestrator-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/alien/.local/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/alien</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
